# systempub.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/SystemPub
Restart=on-failure

[Install]
WantedBy=default.target
